# crontab
# put into /etc/cron.d/configsys

@reboot root /usr/bin/configsys > /tmp/configsys.log

