# crontab
# put into /etc/cron.d/configsys

@reboot root /home/pi/configsys > /tmp/configsys.log

